# simpleposprint.service
[Unit]
Description=SimplePOSPrint Thermal Printer Bridge
After=network.target

[Service]
Type=simple
# Change this to your deployment path (default: /opt/SimplePOSPrint)
WorkingDirectory=/opt/SimplePOSPrint
# Prefer virtualenv if present, else fallback to system python3
ExecStart=/bin/bash -c '[ -d venv ] && source venv/bin/activate && exec python3 spp.py || exec python3 spp.py'
Restart=on-failure
User=root
Group=lp
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
