# simpleposprint.service
[Unit]
Description=SimplePOSPrint Thermal Printer Bridge
After=network.target

[Service]
Type=simple
# Change this to your deployment path (default: /opt/SimplePOSPrint)
WorkingDirectory=/opt/spp
# Prefer virtualenv if present, else fallback to system python3
ExecStart=/opt/spp/venv/bin/python3 /opt/spp/spp.py
Restart=on-failure
User=spp
Group=spp
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
